# etc/services.d/elomapper/finish
#!/usr/bin/execlineb -S0

# Stops all services when "PlayerElo" stops.
s6-svscanctl -t /var/run/s6/services