# etc/services.d/elomapper/run
#!/usr/bin/with-contenv bash
exec elomapper