# etc/services.d/elomapper/run
#!/usr/bin/with-contenv bash
cp -r /views .
cp -r /static .
exec elomapper